# puppet/modules/java/files/etc/profile.d/jdk.sh
# !!! generated by puppet !!!

export JAVA_HOME=/usr/java/default
export JDK_HOME=$JAVA_HOME
export JRE_HOME=$JAVA_HOME/jre
