# puppet/modules/tomcat/templates/etc/init.d/tomcat.erb
#!/bin/sh
#
# !!! generated by puppet !!!

### BEGIN INIT INFO
# Provides:       tomcat
# Required-Start: network 
# Required-Stop:  network
# Default-Start:  345 
# Default-Stop:
# Description:    Tomcat
### END INIT INFO


if echo '\c' | grep -s c >/dev/null 2>&1 ; then
    ECHO_N="echo -n"
    ECHO_C=""
else
    ECHO_N="echo"
    ECHO_C='\c'
fi

TOMCATUSER="<%= tcuser %>"
BASEPATH="<%= tcpath %>"
SCRIPT_HOME=$BASEPATH/bin
if [ -f $BASEPATH/.profile ] ; then
	source $BASEPATH/.profile
fi

case $1 in
	start)
		# Workarount OpenVZ-NFS-Mount Problem
		$ECHO_N "remount all /etc/fstab mounts (mount -a)"
		mount -av | logger
		$ECHO_N "Starting Tomcat ... "
		su - $TOMCATUSER -c "$SCRIPT_HOME/startup.sh"
		echo "ok"
		;;
 	stop)
		$ECHO_N "Stopping Tomcat ... "
		su - $TOMCATUSER -c "$SCRIPT_HOME/shutdown.sh"
		echo "ok"
		;;
 	restart)
  		echo -n "Restarting Tomcat: "
  		$0 stop
  		$0 start
  	;;
 	*)
  	# Print help
  	echo "Usage: $0 {start|stop|restart}" 1>&2
  	exit 1
  	;;
esac
exit 0
